# docker/back/crontab.cron
0 2 * * * nobody cd /app/phonenumberinfo/ && python manage.py update_registry
